# services/psi_exporter.service
[Unit]
Description=PSI Exporter

[Service]
User=root
ExecStart=/usr/local/bin/psi_exporter

[Install]
WantedBy=multi-user.target
